# infra/systemd/database/ogfi-auth-throttle-cleanup.service
[Unit]
Description=OGFI ERP bounded authentication-throttle cleanup
Documentation=file:/opt/ogfi/current/infra/systemd/database/README.md
After=network-online.target
Wants=network-online.target
OnFailure=ogfi-auth-throttle-alert@%n.service

[Service]
Type=oneshot
User=ogfi-runtime
Group=ogfi-runtime
WorkingDirectory=/opt/ogfi/current
EnvironmentFile=/srv/ogfi/config/production.env
RuntimeDirectory=ogfi-auth-throttle-cleanup
RuntimeDirectoryMode=0750
UMask=0027
TimeoutStartSec=55s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/ogfi-auth-throttle-cleanup
ExecStart=/usr/bin/flock --nonblock /run/ogfi-auth-throttle-cleanup/host.lock /usr/bin/pnpm auth-throttle:cleanup
