# infra/systemd/database/ogfi-authorization-denial-health-alert@.service
[Unit]
Description=Journal OGFI ERP authorization-denial monitor failure for %i

[Service]
Type=oneshot
User=ogfi-runtime
Group=ogfi-runtime
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ExecStart=/usr/bin/systemd-cat --identifier=ogfi-authorization-denial-health /usr/bin/echo event=authorization_denial_monitor_failure status=CRITICAL failed_unit=%i hosted_alert_delivery=PENDING
